# nodejs/supervisor/dropbear.sv.conf
[program:dropbear]
command=/usr/sbin/dropbear -F
redirect_stderr=true
