# nodejs/supervisor/node.sv.conf
[program:node]
command=/nodejs/bin/npm start
redirect_stderr=true
